# var/lib/lernstick-exam-client/persistent/etc/systemd/system/keylogger.service
[Unit]
Description=ProXam Keylogger (on-demand via SSE event)

[Service]
Type=simple
ExecStart=/usr/bin/keylogger
ExecStop=/bin/bash -c 'kill $1; tail --pid=$1 -f /dev/null' sh $MAINPID
ExecStopPost=/usr/bin/launch keylogger
Restart=on-failure
RestartSec=5s

# Not enabled by default — started via agent event "startKeylogger"
[Install]
WantedBy=graphical.target
