# var/lib/lernstick-exam-client/persistent/etc/systemd/system/live_stream.service
[Unit]
Description=ProXam live screen stream (ffmpeg WebSocket)

[Service]
Type=simple
ExecStart=/usr/bin/live_stream

[Install]
WantedBy=graphical.target
