# var/lib/lernstick-exam-client/persistent/etc/systemd/system/screenshot.service
[Unit]
Description=ProXam Screenshot (on-demand via SSE event)

[Service]
Type=simple
ExecStart=/usr/bin/screenshot
ExecStop=/bin/bash -c 'kill $1; tail --pid=$1 -f /dev/null' sh $MAINPID
Restart=on-failure
RestartSec=5s

# Not enabled by default — started via agent event "startScreenshot"
[Install]
WantedBy=graphical.target
